# scripts/systemd/kaare-jing.service
[Unit]
Description=Kåre Inner Voice – Jing (0.6B event scanner)
After=network.target kaare.service

[Service]
UMask=0002
User=kaare
WorkingDirectory=/kaare/services/inner_voices
Environment=PATH=/kaare/services/inner_voices/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/kaare/services/inner_voices
ExecStart=/kaare/services/inner_voices/venv/bin/python3 jing_runner.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
